# overlay/etc/init.d/shedos-brain
#!/sbin/openrc-run

# ShedOS brain daemon. Listens on /run/shedos-brain.sock and serves both
# the TUI (JSON-RPC) and legacy nc -U clients (raw text auto-detected).

name="shedos-brain"
description="ShedOS Claude agent daemon (Unix socket RPC)"
command="/usr/bin/python3"
command_args="/opt/shedos/brain.py"
command_background="yes"
pidfile="/run/shedos-brain.pid"
output_log="/var/log/shedos-brain.log"
error_log="/var/log/shedos-brain.log"
respawn_delay=2
supervisor="supervise-daemon"

depend() {
    need net localmount
    before local
}

start_pre() {
    install -d -m 0700 /var/lib/shedos
    install -d -m 0700 /var/lib/shedos/sessions
    # Wait briefly for python3 + py3-* deps installed by alpine-base post.
    # Daemon will fail-fast if /etc/shedos/token is missing — that's OK,
    # the inittab's run-tui.sh prompts the user and writes it, then a
    # supervise-daemon respawn picks it up.
    return 0
}
